# RaspBerry3Bplus/etc/cron.d/crontab
SHELL=/bin/bash

@hourly /home/pi/freenom.sh -u akbal.tk -s home

#* * * * *   /usr/bin/env > /home/pi/cron-env

@reboot /home/pi/scripts/party.cron
